# deploy/libretranslate.service
# /etc/systemd/system/libretranslate.service
[Unit]
Description=LibreTranslate service
After=network.target

[Service]
Type=simple
Environment=LT_LOAD_ONLY=en,fr,es,de,pt,it,ar,hi,sw,am,zh
Environment=LT_DISABLE_FILES_UPLOADS=true
User=libretranslate
Group=libretranslate
# Adjust ExecStart if you installed via pipx or virtualenv
ExecStart=/usr/bin/libretranslate
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
